# systemd/argus-smoke.timer
[Unit]
Description=Run Argus smoke test daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
