# systemd/argus-m0-docker-lockdown.service
[Unit]
Description=Argus M0 dual-stack Docker forwarded-ingress lockdown
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/srv/argus/scripts/argus-m0-docker-lockdown --apply --acknowledge-workload-lockdown
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
