# systemd/argus-domain-agent-legacy-rootful.service
[Unit]
Description=Argus typed compatibility agent for legacy-rootful
After=network.target

[Service]
Type=simple
User=oreo
Group=argus-control
WorkingDirectory=/srv/argus
Environment=ARGUS_ROOT=/srv/argus
Environment=ARGUS_RUNTIME=/srv/argus/runtime/argus/m5
Environment=ARGUS_TRUST_DOMAIN=legacy-rootful
Environment=ARGUS_CAPABILITY_KEY_FILE=/etc/argus/domain-keys/legacy-rootful.key
ExecStart=/usr/bin/python3 /srv/argus/scripts/argus_domain_agent.py
Restart=on-failure
RestartSec=3
UMask=0007
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/argus/runtime/argus -/srv/argus/runtime/audit.log /srv/argus/config /srv/argus/control-plane/dashboard/public /srv/argus/runtime/backups

[Install]
WantedBy=multi-user.target
