# systemd/oreo-smoke.service
[Unit]
Description=Oreo Cloud scheduled smoke test
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/oreo-cloud
ExecStart=/srv/oreo-cloud/scripts/oreo-smoke-scheduled
User=oreo
Group=oreo
