# systemd/oreo-workload-discover.service
[Unit]
Description=Oreo Cloud scheduled workload discovery
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/oreo-cloud
ExecStart=/srv/oreo-cloud/scripts/oreo-workload-discover
User=oreo
Group=oreo
