# templates/systemd/oreo-smoke.timer
[Unit]
Description=Run Oreo Cloud smoke test daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
